# gabbwireless/rootfs/etc/services.d/gabbwireless/run
#!/command/with-contenv bashio
# shellcheck shell=bash disable=SC2034
# =============================================================================
# gabbwireless run script for s6-init
#
# This script automatically detects if it is running as the Home Assistant
# addon or a standard docker environment and sets configuration variables as
# appropriate.
# ==============================================================================

set -euo pipefail

declare DEBUG RUNMODE
declare GABB_USERNAME GABB_PASSWORD REFRESH_RATE
declare MQTT_BROKER MQTT_PORT MQTT_USERNAME MQTT_PASSWORD
declare ENV_VARS

# Delay to keep logs messages from overlapping with s6 logs
sleep 1

# If HASSIO_TOKEN variable exist we are running as addon
if [ -v HASSIO_TOKEN ]; then
  declare -r CONFIG_PATH=/data/options.json

  RUNMODE="addon"
  readonly RUNMODE
  DEBUG=$(bashio::config "debug" || false)
  readonly DEBUG

  # Export MQTT service discovery data for publisher script
  if bashio::services.available 'mqtt'; then
    MQTT_BROKER=$(bashio::services mqtt "host")
    MQTT_PORT=$(bashio::services mqtt "port")
    MQTT_USERNAME=$(bashio::services mqtt "username")
    MQTT_PASSWORD=$(bashio::services mqtt "password")
  fi

  ENV_VARS=$(bashio::config "env_vars")
  if [ -z "$ENV_VARS" ]; then
    bashio::log.info "No additional environment variables found"
  else
    bashio::log.info "Extracted variables ${ENV_VARS}"
  fi

  export ADDONHOSTNAME=$HOSTNAME
else
  RUNMODE="docker"
  readonly RUNMODE

  # Configurable Variables from Environment
  GABB_USERNAME="${GABB_USERNAME:-default_username}"
  GABB_PASSWORD="${GABB_PASSWORD:-default_password}"
  REFRESH_RATE="${REFRESH_RATE:-1}"

  MQTT_BROKER="${MQTT_BROKER:-mqtt.example.com}"
  MQTT_PORT="${MQTT_PORT:-1883}"
  MQTT_USERNAME="${MQTT_USERNAME:-mqtt_user}"
  MQTT_PASSWORD="${MQTT_PASSWORD:-mqtt_password}"
fi

readonly GABB_USERNAME GABB_PASSWORD REFRESH_RATE
readonly MQTT_BROKER MQTT_PORT MQTT_PASSWORD
readonly ENV_VARS

cd /app/

echo "-------------------------------------------------------"
echo "gabbwireless version: $(awk '/__version__/ { print $3}' gabb/__init__.py | sed 's/"//g')"
echo "Python version $(python3 --version | awk '{print $2}')"
echo "-------------------------------------------------------"

bashio::log.info "Starting Gabb Wireless..."
if [ -z "$ENV_VARS" ]; then
  bashio::log.info "Running Gabb Wireless env_vars"
  exec python3 ./gabb_mqtt_publisher.py
else
  env "${ENV_VARS[@]}" python3 ./gabb_mqtt_publisher.py
fi
